# debian/fast-api-app/DEBIAN/prerm
#!/bin/bash

set -e

systemctl stop fastapi-app.service || true

systemctl disable fastapi-app.service || true

rm -f /etc/systemd/system/fastapi-app.service

